# collectd/supervisord.conf
[supervisord]
nodaemon=true

[program:collectd]
command=/usr/sbin/collectd -C /etc/collectd/collectd.conf -f
numprocs=1
autostart=true
autorestart=true
