# influxdb/supervisord.conf
[supervisord]
nodaemon=true

[program:influxdb]
command=/usr/bin/influxdb -config /opt/influxdb/current/config.toml
numprocs=1
autostart=true
autorestart=true